# modules/storm_supervisord/templates/storm-ui.ini.erb
[program:storm-ui]
command=env PATH=$PATH:/bin:/usr/bin/:/usr/local/jdk1.7.0_79/bin/ JAVA_HOME=/usr/local/jdk1.7.0_79/ storm ui
directory=/home/storm
autorestart=true
user=storm